# services/cp-soc-dht11.service
[Unit]
Description=CP-SOC DHT11 Sensor Logger
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/cyber-physical-soc
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/iot-venv/bin/python /home/pi/cyber-physical-soc/sensors/dht11_to_jsonl.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
